# deploy/systemd/eink-frame.service
[Unit]
Description=eInk Photo Frame Web App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inky
Group=inky
WorkingDirectory=/home/inky/inky-frame
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/inky/inky-frame/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 80

Restart=always
RestartSec=3

# Allow non-root bind to privileged port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
